# bin/elastic-connector.service
[Unit]
Description=Elasticsearch IOTICS connector
After=syslog.target network.target
[Service]
SuccessExitStatus=143
User=ec2-user
Group=ec2-user
Type=forking

ExecStart=/opt/iotics/elastic-sink-1.0-SNAPSHOT/run.sh
ExecStop=/bin/kill -15 $MAINPID

[Install]
WantedBy=multi-user.target

